# supervisord.conf
[program:apache2]
command=/bin/bash -c ". /etc/apache2/envvars && /usr/sbin/apache2 -D FOREGROUND"
numprocs=1
autostart=true
autorestart=true

[program:mysqld]
command=/usr/bin/mysqld_safe
numprocs=1
autostart=true
autorestart=true
